# roles/canary/templates/digitalocean_api_canary.timer.j2
[Unit]
Description=Check DigitalOcean API connectivity

[Install]
WantedBy=timers.target

[Timer]
OnCalendar=*:*
AccuracySec=1us
